# pkg/osx/osx-pkg/Contents/MacOS/Library/LaunchDaemons/com.qosient.rastream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>KeepAlive</key>
                <true/>
        <key>Program</key>
                <string>/usr/local/bin/rastream</string>
        <key>Label</key>
                <string>com.qosient.rastream</string>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/local/bin/rastream</string>
                <string>-S</string>
                <string>localhost:562</string>
                <string>-B</string>
                <string>15s</string>
                <string>-M</string>
                <string>time</string>
                <string>5m</string>
                <string>-f</string>
                <string>/usr/local/bin/rastream.sh</string>
                <string>-w</string>
                <string>/usr/local/argus/archive/$srcid/%Y/%m/%d/argus.%Y.%m.%d.%H.%M.%S</string>
        </array>
        <key>Disabled</key>
                <false/>
</dict>
</plist>
