# pkg/rhel/systemd/behave-nightly.service
[Unit]
Description=Process Behavioral Alarms and Alerts

[Service]
Type=oneshot
ExecStart=/sbin/behave-nightly
TimeoutStartSec=1hr

