# pkg/rhel/systemd/dbaseRollup-nightly.service
[Unit]
Description=Rollup Argus Database Tables By Month and Year

[Service]
Type=oneshot
ExecStart=/sbin/dbaseRollup-nightly
TimeoutStartSec=1hr

