# pkg/rhel/systemd/hostsData-nightly.service
[Unit]
Description=Process Matrix Database tables

[Service]
Type=oneshot
ExecStart=/sbin/matrixData-nightly
TimeoutStartSec=1hr

