# pkg/rhel/systemd/radhcp-nightly.service
[Unit]
Description=Process DHCP transactions

[Service]
Type=oneshot
ExecStart=/sbin/radhcp-nightly
TimeoutStartSec=1hr

