# pkg/rhel/systemd/radns-nightly.service
[Unit]
Description=Process DNS transactions

[Service]
Type=oneshot
ExecStart=/sbin/radns-nightly
TimeoutStartSec=1hr

