# pkg/rhel/systemd/rahisto-delsql.service
[Unit]
Description=Remove old histograms from database

[Service]
Type=oneshot
ExecStart=/usr/bin/rahisto-delsql -t -6M
TimeoutStartSec=1hr
