# pkg/rhel/systemd/rahosts-nightly.service
[Unit]
Description=Process Hosts Daily COI data

[Service]
Type=oneshot
ExecStart=/usr/bin/rahostsdaily
TimeoutStartSec=1hr

