# pkg/rhel/systemd/raports-nightly.service
[Unit]
Description=Process Port Utilization Data

[Service]
Type=oneshot
ExecStart=/usr/bin/raportsdaily
TimeoutStartSec=1hr

