# pkg/rhel/systemd/rarp-nightly.service
[Unit]
Description=Process ARP transactions

[Service]
Type=oneshot
ExecStart=/sbin/rarp-nightly
TimeoutStartSec=1hr

