# pkg/osx/Library/LaunchDaemons/com.qosient.rainventory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
   <dict>
      <key>KeepAlive</key>
         <true/>
      <key>Label</key>
         <string>com.qosient.rainventory</string>
      <key>ProgramArguments</key>
         <array>
            <string>/usr/local/bin/rasqlinsert</string>
            <string>-M</string>
            <string>cache</string>
            <string>rmon</string>
            <string>time</string>
            <string>1d</string>
            <string>-m</string>
            <string>saddr</string>
            <string>smac</string>
            <string>sid</string>
            <string>inf</string>
            <string>-s</string>
            <string>stime</string>
            <string>ltime</string>
            <string>dur</string>
            <string>sid</string>
            <string>inf</string>
            <string>smac</string>
            <string>saddr</string>
            <string>pkts</string>
            <string>bytes</string>
            <string>pcr</string>
            <string>-S</string>
            <string>localhost</string>
            <string>-w</string>
            <string>mysql://root@localhost/inventory/ipAddrs_%Y_%m_%d</string>
            <string>-</string>
            <string>ip</string>
         </array>
   </dict>
</plist>
